# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gateway]
command=dotnet FlowBoard.Gateway.dll --urls "http://0.0.0.0:10000"
directory=/app/gateway
autostart=true
autorestart=true
startsecs=1
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth]
command=dotnet FlowBoard.Auth.dll --urls "http://0.0.0.0:5048"
directory=/app/auth
autostart=true
autorestart=true
startsecs=1
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:workspace]
command=sh -c "sleep 5 && dotnet FlowBoard.Workspace.dll --urls 'http://0.0.0.0:5277'"
directory=/app/workspace
autostart=true
autorestart=false
startsecs=6
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:board]
command=sh -c "sleep 10 && dotnet FlowBoard.Board.dll --urls 'http://0.0.0.0:5030'"
directory=/app/board
autostart=true
autorestart=false
startsecs=11
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:list]
command=sh -c "sleep 15 && dotnet FlowBoard.List.dll --urls 'http://0.0.0.0:5266'"
directory=/app/list
autostart=true
autorestart=false
startsecs=16
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:card]
command=sh -c "sleep 20 && dotnet FlowBoard.Card.dll --urls 'http://0.0.0.0:5024'"
directory=/app/card
autostart=true
autorestart=false
startsecs=21
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:comment]
command=sh -c "sleep 25 && dotnet FlowBoard.Comment.dll --urls 'http://0.0.0.0:5198'"
directory=/app/comment
autostart=true
autorestart=false
startsecs=26
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:label]
command=sh -c "sleep 30 && dotnet FlowBoard.Label.dll --urls 'http://0.0.0.0:5027'"
directory=/app/label
autostart=true
autorestart=false
startsecs=31
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


